# infra/systemd/fasttravel-snapshot.timer
# ============================================================================
# Fires fasttravel-snapshot.service at 06:00 and 18:00 Europe/Kyiv.
#
# Note: the host timezone is set to Europe/Kyiv in cloud-init.yml, so plain
# OnCalendar uses local Kyiv time. The `Timezone=` field on OnCalendar requires
# systemd ≥ 252 and Ubuntu 22.04 ships systemd 249 — DO NOT use it here.
# ============================================================================
[Unit]
Description=Twice-daily FastTravel price snapshot
Requires=fasttravel-snapshot.service

[Timer]
OnCalendar=*-*-* 06,18:00:00
Persistent=true
AccuracySec=1min
Unit=fasttravel-snapshot.service

[Install]
WantedBy=timers.target
