# infra/systemd/fasttravel-snapshot.service
# ============================================================================
# FastTravel — price snapshot worker (one-shot).
# Triggered twice daily by fasttravel-snapshot.timer.
# ============================================================================
[Unit]
Description=FastTravel price snapshot job
After=fasttravel-stack.service network-online.target
Wants=network-online.target
Requires=fasttravel-stack.service

[Service]
Type=oneshot
WorkingDirectory=/opt/fasttravel
EnvironmentFile=-/opt/fasttravel/.env
# Runs the snapshot job inside the scheduler container.
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml exec -T scheduler python -m apps.scheduler.src.jobs.snapshot_hot
# Snapshot of 300 hotels × 90 days × 3 nights × 2 meals × 3 operators with
# concurrency=5 takes ~20 min. Give it 45 min before SIGTERM.
TimeoutStartSec=2700
StandardOutput=journal
StandardError=journal
